# CentOS_6/vsftpd/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
/usr/sbin/chroot /var/ftp/ftpjail /usr/sbin/vsftpd
/usr/local/bin/snort -D -c /etc/snort/snort.conf -i eth0
